# server/shared-server-deploy/systemd/gaozhikai-us-stock-crash-openai-refresh-0755.timer
[Unit]
Description=07:55 refresh for 高志凯美股崩盘概率预测程序（openai版）

[Timer]
OnCalendar=*-*-* 07:55:00 Asia/Shanghai
Unit=gaozhikai-us-stock-crash-openai-refresh@0755.service
Persistent=true

[Install]
WantedBy=timers.target
