# server/shared-server-deploy/systemd/gaozhikai-us-stock-crash-openai-refresh-0745.timer
[Unit]
Description=07:45 refresh for 高志凯美股崩盘概率预测程序（openai版）

[Timer]
OnCalendar=*-*-* 07:45:00 Asia/Shanghai
Unit=gaozhikai-us-stock-crash-openai-refresh@0745.service
Persistent=true

[Install]
WantedBy=timers.target
